# misc/pokemongbackend.service
[Unit]
Description=Pokemong's backend

[Service]
Type=simple
ExecStart=/usr/bin/pokemongbackend
